# packaging/osx/postinstall
#!/bin/bash
mkdir -p /var/log/snclient
launchctl load /Library/LaunchDaemons/com.snclient.snclient.plist
