# packaging/debian/postinst
#!/bin/bash
set -e

if [[ `id snclient -u -n` != "snclient" ]] ; then
  adduser --system --home /nonexistent --no-create-home snclient
fi

if [[ ! -d /var/log/snclient ]]; then
  mkdir /var/log/snclient
fi

chown snclient /etc/snclient/*
chown snclient:adm /var/log/snclient

case "$1" in
    configure|abort-upgrade)
      if [ -x "/usr/bin/deb-systemd-helper" ]; then
        deb-systemd-helper unmask snclient.service
        if deb-systemd-helper --quiet was-enabled snclient.service; then
          deb-systemd-helper enable snclient.service
        else
          deb-systemd-helper update-state snclient.service
        fi
      fi

      systemctl --system daemon-reload
      deb-systemd-invoke start snclient
    ;;

    abort-remove|abort-deconfigure)
    ;;

    *)
        echo "postinst called with unknown argument \`$1'" >&2
        exit 1
    ;;
esac

# dh_installdeb will replace this with shell code automatically
# generated by other debhelper scripts.

#DEBHELPER#

exit 0