# packaging/debian/postrm
#!/bin/sh
set -e
if [ -x "/usr/bin/systemctl" ] ; then
	systemctl --system daemon-reload >/dev/null 2>&1 || true
fi

if [ "$1" = "purge" ] ; then
	update-rc.d snclient remove
fi
