# packaging/snclient.post-install
#!/bin/sh

rc-update -q add snclient default
rc-service -q snclient start
