# packaging/snclient.service
[Unit]
Description=SNClient Agent
After=network.target

[Service]
Type=simple
User=snclient
Group=snclient

AmbientCapabilities=CAP_SETUID CAP_SETGID

CacheDirectory=snclient
CacheDirectoryMode=0750

# Last-resort protection for the complete service cgroup, including child
# processes. The in-process 512 MiB RSS watcher remains the primary safeguard.
MemoryAccounting=yes
MemoryHigh=1536M
MemoryMax=2G

# Treat a cgroup OOM as a failure of the complete service and restart cleanly.
OOMPolicy=kill
KillMode=control-group

Restart=on-failure
RestartSec=10

WorkingDirectory=/var/lib/snclient
ExecStart=/usr/bin/snclient --config=/etc/snclient/snclient.ini
ExecReload=/bin/kill -HUP $MAINPID

[Install]
WantedBy=multi-user.target
