# packaging/snclient.service
[Unit]
Description=SNClient+ Agent
After=network.target

[Service]
Type=simple
User=root

Restart=on-failure
RestartSec=10
startLimitIntervalSec=60

WorkingDirectory=/
ExecStart=/usr/bin/snclient --config=/etc/snclient/snclient.ini
ExecStartPre=/bin/mkdir -p /var/log/snclient

[Install]
WantedBy=multi-user.target
